# templates/graphite-api.service.j2
[Unit]
Description=graphite-api

[Service]
ExecStart=/data/graphite-api/bin/uwsgi --ini /data/etc/graphite.ini
