# ansible/uwsgi.service.j2
[Unit]
Description=uWSGI

[Service]
ExecStart=/data/venv/bin/uwsgi --ini /data/etc/uwsgi.ini
